# confidential-voting/infra/vps/fhevm-frontend.service
[Unit]
Description=FHEVM Vote - Frontend server (no-cache + reverse proxy to Hardhat)
After=fhevm-hardhat.service fhevm-relayer.service
BindsTo=fhevm-hardhat.service

[Service]
Type=simple
User=fhevm
WorkingDirectory=/opt/fhevm-vote/confidential-voting
ExecStart=/usr/bin/python3 frontend_server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
